# server/home/godot/.config/systemd/user/godot-benchmarks.service
[Service]
ExecStart=/home/godot/godot-benchmarks/run-benchmarks.sh

[Install]
WantedBy=multi-user.target
