# server/home/godot/.config/systemd/user/godot-benchmarks.timer
[Unit]
Description=Run godot-benchmarks

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
