# server/home/godot/.config/systemd/user/xorg.service
[Service]
ExecStart=xinit
Restart=always

[Install]
WantedBy=default.target
