# server/home/godot/.config/systemd/user/xorg.service
[Service]
ExecStart=xinit

[Install]
WantedBy=multi-user.target
